# board/eeupdate/rootfs_overlay/etc/profile.d/serial-tty-resize.sh
# Fix the terminal size on serial consoles
if tty | grep -q ^/dev/ttyS[0-9] ; then
	eval $(resize)
fi
